# dot_config/profile.d/000.keep.sh
# Note: profile refers to something that runs only in login shells (spawned by your terminal, not just a shell script,
# for example).
