# dot_config/profile.d/900.cheat-sheet.sh
# Cheat sheet
if test -r "$HOME/bash-cheat-sheet.txt"; then
  \cat "$HOME/bash-cheat-sheet.txt"
else
  echo "'$HOME/bash-cheat-sheet.txt' not available"
fi
