# .config/profile.d/300.bash-fixup-stty.sh
### Bash setup

# https://unix.stackexchange.com/questions/12107/how-to-unfreeze-after-accidentally-pressing-ctrl-s-in-a-terminal
stty -ixon -ixoff
