# .config/profile.d/301.bash-git-completion.sh
# Bash completion
export BASH_COMPLETION_COMPAT_DIR="/usr/local/etc/bash_completion.d"
[[ -r "/usr/local/etc/profile.d/bash_completion.sh" ]] && . "/usr/local/etc/profile.d/bash_completion.sh"

# Git completion
# shellcheck disable=SC1090
source "$HOME/git-completion.bash"
