# .config/profile.d/303.bash-autojump.sh
# Autojump!
if test -r /usr/local/etc/profile.d/autojump.sh; then
    source /usr/local/etc/profile.d/autojump.sh
fi
