# .config/profile.d/307.bash-travis.sh
# added by travis gem
if test -f "$HOME/.travis/travis.sh"; then
  # shellcheck disable=SC1090
  source "$HOME/.travis/travis.sh"
fi
