# .config/profile.d/404.func.termbin.sh
function tb() {
  echo "Sending '$*' to termbin..."
  nc termbin.com 9999 < "$*"
}
